# docs/ImpftermineChecker.service
[Unit]
Description=Checkt auf Impftermine

[Service]
Type=simple
ExecStart=impftermin-telegram

[Install]
WantedBy=basic.target
