# docs/ImpftermineChecker.timer
[Unit]
Description=Impfchecker Timer

[Timer]
OnCalendar=minutely
Persistent=true
Unit=ImpftermineChecker.service

[Install]
WantedBy=basic.target
